# config/puma-auto-glossary.service
# ============================================================================
# SYSTEMD SERVICE FILE FOR PUMA WEB SERVER
# ============================================================================
# This file defines the systemd service that manages Puma in production.
#
# INSTALLATION:
# sudo cp /opt/auto-glossary/current/config/puma.service /etc/systemd/system/puma-auto-glossary.service
# sudo systemctl daemon-reload
# sudo systemctl enable puma-auto-glossary.service
#
# MANAGEMENT COMMANDS:
# - Start:   sudo systemctl start puma-auto-glossary.service
# - Stop:    sudo systemctl stop puma-auto-glossary.service
# - Restart: sudo systemctl restart puma-auto-glossary.service
# - Status:  sudo systemctl status puma-auto-glossary.service
# - Logs:    sudo journalctl -xeu puma-auto-glossary.service
#
# AUTOMATIC DEPLOYMENT:
# Capistrano automatically restarts Puma after deployment via:
# config/deploy.rb: after "deploy:published", "systemd_puma:restart"
#
# ⚠️  CRITICAL: Changes to this file require:
# 1. Copy to /etc/systemd/system/puma-auto-glossary.service
# 2. sudo systemctl daemon-reload
# 3. sudo systemctl restart puma-auto-glossary.service
# ============================================================================

[Unit]
Description=Puma HTTP Server for auto-glossary (production)
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5
OnFailure=notify-on-failure@%n.service

[Service]
Type=simple
User=deploy
MemoryMax=512M
CPUQuota=80%
WorkingDirectory=/opt/auto-glossary/current

# ============================================================================
# ENVIRONMENT VARIABLES - CRITICAL FOR RAILS OPERATION
# ============================================================================
# RAILS_ENV: MUST be "production" - determines which config files are loaded
# RBENV_ROOT: Required for rbenv to find Ruby installation
# PATH: Ensures rbenv shims are used (Ruby 3.4.3 as configured)
#
# ⚠️  MISSING ENVIRONMENT VARIABLES:
# The following are loaded from /opt/auto-glossary/shared/.env via Capistrano:
# - MYSQL_DATABASE, MYSQL_USER, MYSQL_PASSWORD: Database credentials
# - RAILS_MASTER_KEY: Decrypts credentials.yml.enc
#
# TO ADD ENVIRONMENT VARIABLES:
# 1. Edit: /opt/auto-glossary/shared/.env (NOT this file)
# 2. Format: KEY=value (one per line, no export statement)
# 3. Restart: sudo systemctl restart puma-auto-glossary.service
#
# NOTE: .env is loaded by Rails, not systemd. See config/deploy.rb.
# ============================================================================
Environment=RAILS_ENV=production
Environment=RBENV_ROOT=/home/deploy/.rbenv
Environment=PATH=/home/deploy/.rbenv/shims:/home/deploy/.rbenv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# ============================================================================
# PROCESS MANAGEMENT
# ============================================================================
# ExecStartPre: Cleans up stale socket/pid files before starting
# ExecStart: Starts Puma using bundler to ensure correct gem versions
# ExecReload: Graceful restart (USR1 signal) - workers finish current requests
# ExecStop: Graceful shutdown (SIGTERM) - allows cleanup before exit
#
# ⚠️  IMPORTANT: Puma configuration is in /opt/auto-glossary/current/config/puma.rb
# Do NOT add puma CLI flags here - all config should be in puma.rb
# ============================================================================
ExecStartPre=/bin/rm -f /opt/auto-glossary/shared/tmp/sockets/puma.sock
ExecStartPre=/bin/rm -f /opt/auto-glossary/shared/tmp/pids/puma.pid
ExecStartPre=/bin/rm -f /opt/auto-glossary/shared/tmp/pids/puma.state
ExecStart=/home/deploy/.rbenv/shims/bundle exec puma -C /opt/auto-glossary/current/config/puma.rb
ExecReload=/bin/kill -USR1 $MAINPID
ExecStop=/bin/kill -SIGTERM $MAINPID

# ============================================================================
# RESTART POLICY - ENSURES HIGH AVAILABILITY
# ============================================================================
# Restart=always: Automatically restart if Puma crashes
# RestartSec=10: Wait 10 seconds before restarting (prevents rapid restart loops)
# KillMode=mixed: Send SIGTERM to main process, SIGKILL to remaining after timeout
# TimeoutStopSec=30: Wait 30 seconds for graceful shutdown before force kill
# StartLimitBurst=5: Allow 5 restart attempts within StartLimitIntervalSec
# StartLimitIntervalSec=300: Reset restart counter after 5 minutes
#
# This configuration prevents 500 errors from persisting after crashes
# ============================================================================
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# ============================================================================
# LOGGING
# ============================================================================
# StandardOutput: Puma stdout (info messages, requests if configured)
# StandardError: Puma stderr (errors, exceptions, backtraces)
# SyslogIdentifier: Tag for journalctl filtering
#
# VIEW LOGS:
# - Systemd: sudo journalctl -xeu puma-auto-glossary.service
# - Stdout:  tail -f /opt/auto-glossary/shared/log/puma_stdout.log
# - Stderr:  tail -f /opt/auto-glossary/shared/log/puma_stderr.log
# - Rails:   tail -f /opt/auto-glossary/shared/log/production.log (NOT managed here)
# ============================================================================
StandardOutput=append:/opt/auto-glossary/shared/log/puma_stdout.log
StandardError=append:/opt/auto-glossary/shared/log/puma_stderr.log
SyslogIdentifier=puma-auto-glossary

[Install]
WantedBy=multi-user.target
